# systemd/phase1/picoclaw-phase1-recovery.service
[Unit]
Description=PicoClaw Phase1 Missed Update Recovery
Wants=network-online.target
After=network-online.target picoclaw.service

[Service]
Type=oneshot
User=root
EnvironmentFile=-/etc/picoclaw/phase1.env
ExecStart=/usr/bin/python3 /root/picoclaw/scripts/phase1/missed_update_recovery.py

[Install]
WantedBy=multi-user.target
